# backend/data/corpus/techniques.json
[
  {
    "id": "T1110",
    "title": "Brute Force",
    "tactic": "credential-access",
    "mitigations": "Enforce account lockout policies. Require multi-factor authentication (MFA) on all accounts. Monitor and alert on repeated failed authentication attempts. Use CAPTCHA to slow automated attacks. Restrict access from unusual geographies or IP ranges."
  },
  {
    "id": "T1078",
    "title": "Valid Accounts",
    "tactic": "defense-evasion",
    "mitigations": "Audit account usage and privileges regularly. Enforce least-privilege principles. Monitor for anomalous login times and locations. Implement privileged access management (PAM). Rotate credentials on compromise suspicion. Disable inactive accounts promptly."
  },
  {
    "id": "T1566",
    "title": "Phishing",
    "tactic": "initial-access",
    "mitigations": "Deploy email filtering and anti-phishing controls. Train users to identify suspicious emails. Use DMARC, DKIM, SPF to authenticate email senders. Sandbox email attachments before delivery. Block known malicious domains at DNS and proxy layers."
  },
  {
    "id": "T1059",
    "title": "Command and Scripting Interpreter",
    "tactic": "execution",
    "mitigations": "Restrict script execution via application allow-listing. Disable or limit interpreter access where not required. Log and monitor interpreter invocations. Use endpoint detection to alert on suspicious script usage."
  },
  {
    "id": "T1059.001",
    "title": "PowerShell",
    "tactic": "execution",
    "mitigations": "Enforce PowerShell Constrained Language Mode. Enable script block logging and module logging. Disable interactive PowerShell for standard users. Alert on encoded commands (-enc/-EncodedCommand) and base64 payloads. Use AMSI for runtime inspection."
  },
  {
    "id": "T1021",
    "title": "Remote Services",
    "tactic": "lateral-movement",
    "mitigations": "Limit remote service exposure to the internet. Require MFA for remote access. Audit remote access logs. Use VPN with strong authentication. Apply network segmentation to isolate sensitive systems."
  },
  {
    "id": "T1486",
    "title": "Data Encrypted for Impact",
    "tactic": "impact",
    "mitigations": "Maintain offline and offsite backups. Test backups regularly. Deploy endpoint protection with ransomware detection. Use file integrity monitoring. Isolate infected hosts immediately. Have an incident response plan ready."
  },
  {
    "id": "T1190",
    "title": "Exploit Public-Facing Application",
    "tactic": "initial-access",
    "mitigations": "Keep public-facing applications patched and up-to-date. Deploy a web application firewall (WAF). Perform regular vulnerability scanning. Reduce the attack surface by disabling unused services. Use intrusion detection on perimeter traffic."
  },
  {
    "id": "T1053",
    "title": "Scheduled Task/Job",
    "tactic": "persistence",
    "mitigations": "Audit scheduled tasks and cron jobs regularly. Monitor for new or modified scheduled tasks. Restrict permissions to create scheduled tasks. Alert on scheduled task creation by non-privileged users."
  },
  {
    "id": "T1046",
    "title": "Network Service Discovery",
    "tactic": "discovery",
    "mitigations": "Monitor for port-scanning and host-enumeration patterns (e.g. nmap, masscan). Segment the network so discovery blast radius is limited. Alert on internal hosts scanning critical subnets. Correlate scans with subsequent lateral movement attempts."
  },
  {
    "id": "T1071",
    "title": "Application Layer Protocol",
    "tactic": "command-and-control",
    "mitigations": "Inspect outbound TLS traffic for beaconing patterns (regular intervals, small payloads). Block known C2 domains and IPs at DNS, proxy, and firewall layers. Use threat-intel feeds to flag outbound connections to malicious infrastructure. Hunt for DGAs and look-alike domains."
  },
  {
    "id": "T1048",
    "title": "Exfiltration Over Alternative Protocol",
    "tactic": "exfiltration",
    "mitigations": "Monitor outbound data volume per host and per destination. Alert on large or sustained transfers to unknown/unregistered ASNs. Block uncommon protocols on perimeter egress. DLP on sensitive directories. Investigate after-hours bulk transfers."
  },
  {
    "id": "T1218",
    "title": "System Binary Proxy Execution",
    "tactic": "defense-evasion",
    "mitigations": "Monitor signed binaries (e.g. mshta, rundll32, certutil) used to execute suspicious content. Alert on signed binaries spawned by unusual parents (e.g. office apps spawning interpreters). Use behavioral EDR rules. Restrict binary execution on endpoints."
  },
  {
    "id": "T1133",
    "title": "External Remote Services",
    "tactic": "initial-access",
    "mitigations": "Require MFA for all external remote access. Monitor for impossible-travel and anomalous geo logins. Limit VPN/RDP exposure via jump hosts. Alert on privileged logins from new geographies. Enforce time-bound access for privileged accounts."
  },
  {
    "id": "T1027",
    "title": "Obfuscated Files or Information",
    "tactic": "defense-evasion",
    "mitigations": "Inspect and decode obfuscated content (base64, XOR, encoded payloads) at the endpoint and email gateway. Alert on encoded command-line invocations. Use AMSI and script-block logging. Block macros that spawn encoded interpreters."
  }
]
